# DEBIAN/postinst
#!/usr/bin/env bash

chmod +x /usr/bin/autoupdater
systemctl enable autoupdater.timer